# work/profile.d/40_acquia_support_bin.sh
export PATH="${HOME}/src/Support-Tools/bin:${PATH}"

