# common/profile.d/99_powerline_prompt.sh
export PROMPT_COMMAND="history -a; history -c; history -r; _powerline_ps1"

